# Stuff/digestbotbackup.service
[Unit]
Description=Backup of DigestBot Stack Log

[Service]
Type=simple
ExecStart=/root/backup_digest_log.sh

[Install]
WantedBy=multi-user.target
